# launchd/com.bathjobs.scraper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.bathjobs.scraper</string>
  <key>ProgramArguments</key>
  <array>
    <string>PROJECT_DIR/venv/bin/python</string>
    <string>PROJECT_DIR/scraper.py</string>
  </array>
  <key>WorkingDirectory</key>
  <string>PROJECT_DIR</string>
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key>   <integer>6</integer>
      <key>Minute</key> <integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key>   <integer>18</integer>
      <key>Minute</key> <integer>0</integer>
    </dict>
  </array>
  <key>StandardOutPath</key>
  <string>PROJECT_DIR/logs/scraper.log</string>
  <key>StandardErrorPath</key>
  <string>PROJECT_DIR/logs/scraper.log</string>
</dict>
</plist>
